# mail/mta/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
process_name = master
directory = /etc/postfix
command = /usr/sbin/postfix -c /etc/postfix start
startsecs = 0
autorestart = false

[program:gross]
command = /usr/sbin/grossd -u gross -d -f /etc/grossd.conf
startsecs = 0

[program:dovecot]
command = /usr/sbin/dovecot
startsecs = 0

[program:spamd]
command = /usr/sbin/spamd --ipv4 --listen 127.0.0.1:783 --allow-tell --allowed-ips=127.0.0.1
startsecs = 0
